# examples/12-RasPi5_i226-grandmaster/etc/systemd/system/chrony.service.d/override.conf
[Service]
AmbientCapabilities=CAP_SYS_TIME CAP_SYS_NICE
CapabilityBoundingSet=CAP_SYS_TIME CAP_SYS_NICE
NoNewPrivileges=no
